# statsd/run_once.supervisord.conf
[program:set_graphite_datasource] 
command=sh /root/set_graphite_datasource.sh
autorestart=false 
