# statsd/run_once.supervisord.conf
[supervisord]
nodaemon=true

[program:set_graphite_datasource] 
command=sh /root/set_graphite_datasource.sh
autostart=false 
autorestart=false 
startretries=0 
